# gost/gost.service.example
[Unit]
Description=Gost Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/tunnel/gost/gost -C /opt/tunnel/gost/config.json
ExecStartPre=/opt/tunnel/pre-start.sh
ExecStopPost=/opt/tunnel/post-stop.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
